# scripts/install/linux/debian/install.sh
#!/bin/bash

# Script de Instala칞칚o para Linux Debian/Ubuntu - Finan칞as UP
# Vers칚o: 1.0.0
# Data: 19/10/2025

set -e  # Parar em caso de erro

# Cores para output
RED='\033[0;31m'
GREEN='\033[0;32m'
YELLOW='\033[1;33m'
BLUE='\033[0;34m'
NC='\033[0m' # No Color

# Fun칞칚o para log
log() {
    echo -e "${BLUE}[INFO]${NC} $1"
}

success() {
    echo -e "${GREEN}[SUCCESS]${NC} $1"
}

error() {
    echo -e "${RED}[ERROR]${NC} $1"
    exit 1
}

warning() {
    echo -e "${YELLOW}[WARNING]${NC} $1"
}

# Verificar se 칠 Debian/Ubuntu
if ! (grep -q "Debian" /etc/os-release || grep -q "Ubuntu" /etc/os-release); then
    error "Este script 칠 espec칤fico para Debian/Ubuntu. Para outras distribui칞칫es, use o script apropriado."
fi

# Verificar se 칠 root ou tem sudo
if [[ $EUID -eq 0 ]]; then
    SUDO_CMD=""
else
    if ! command -v sudo &> /dev/null; then
        error "sudo n칚o est치 instalado ou configurado."
    fi
    SUDO_CMD="sudo"
fi

log "游 Iniciando instala칞칚o das depend칡ncias para Finan칞as UP em Debian/Ubuntu..."

# Atualizar reposit칩rios
log "游닍 Atualizando reposit칩rios..."
$SUDO_CMD apt update || error "Falha ao atualizar reposit칩rios"

# Instalar curl e outros utilit치rios necess치rios
log "游닍 Instalando utilit치rios b치sicos..."
$SUDO_CMD apt install -y curl gnupg2 software-properties-common || error "Falha ao instalar utilit치rios"

# Instalar Node.js (usando NodeSource para vers칚o recente)
log "游댢 Verificando Node.js..."
if ! command -v node &> /dev/null || ! command -v npm &> /dev/null; then
    log "游닍 Instalando Node.js 18+..."
    curl -fsSL https://deb.nodesource.com/setup_18.x | $SUDO_CMD -E bash -
    $SUDO_CMD apt install -y nodejs || error "Falha ao instalar Node.js"
    success "Node.js e npm instalados"
else
    success "Node.js e npm j치 est칚o instalados"
fi

# Instalar Git
log "游댢 Verificando Git..."
if ! command -v git &> /dev/null; then
    log "游닍 Instalando Git..."
    $SUDO_CMD apt install -y git || error "Falha ao instalar Git"
    success "Git instalado"
else
    success "Git j치 est치 instalado"
fi

# Instalar PostgreSQL
log "游댢 Verificando PostgreSQL..."
if ! command -v psql &> /dev/null; then
    log "游닍 Instalando PostgreSQL..."
    $SUDO_CMD apt install -y postgresql postgresql-contrib || error "Falha ao instalar PostgreSQL"

    # Iniciar servi칞o PostgreSQL
    $SUDO_CMD systemctl enable postgresql
    $SUDO_CMD systemctl start postgresql
    success "PostgreSQL instalado e iniciado"
else
    success "PostgreSQL j치 est치 instalado"
fi

# Verificar vers칫es
log "游댌 Verificando vers칫es instaladas..."
echo ""

if command -v node &> /dev/null; then
    echo -e "Node.js: $(node --version)"
else
    error "Node.js n칚o encontrado ap칩s instala칞칚o"
fi

if command -v npm &> /dev/null; then
    echo -e "npm: $(npm --version)"
else
    error "npm n칚o encontrado ap칩s instala칞칚o"
fi

if command -v git &> /dev/null; then
    echo -e "Git: $(git --version)"
else
    error "Git n칚o encontrado ap칩s instala칞칚o"
fi

if command -v psql &> /dev/null; then
    echo -e "PostgreSQL: $(psql --version)"
else
    error "PostgreSQL n칚o encontrado ap칩s instala칞칚o"
fi

# Configurar PostgreSQL (usu치rio padr칚o)
log "游댢 Configurando PostgreSQL..."
$SUDO_CMD -u postgres createuser --interactive --pwprompt $(whoami) 2>/dev/null || warning "Usu치rio PostgreSQL pode j치 existir"
$SUDO_CMD -u postgres createdb $(whoami) 2>/dev/null || warning "Banco de dados pode j치 existir"

success "Instala칞칚o conclu칤da!"

log "游늶 Pr칩ximos passos:"
echo "1. Clone o reposit칩rio:"
echo "   git clone https://github.com/alphasinfo/financas-up.git"
echo "   cd financas-up"
echo ""
echo "2. Configure o ambiente:"
echo "   npm install"
echo ""
echo "3. Configure o banco:"
echo "   npx prisma generate"
echo "   npx prisma db push"
echo ""
echo "4. Execute o projeto:"
echo "   npm run dev"
echo ""
echo "游닀 Consulte a documenta칞칚o em docs/01-INICIO-RAPIDO.md para mais detalhes"

# Finalizar
success "游꿀 Ambiente Debian/Ubuntu configurado para Finan칞as UP!"
